# edge-device-stack/edge-agent/systemd/edge-agent.service.tmpl
[Unit]
Description=KIN-GOLIATH Edge Agent
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=@SERVICE_USER@
Group=@SERVICE_GROUP@
WorkingDirectory=@EDGE_DIR@
EnvironmentFile=@ENV_FILE@
ExecStart=@UVICORN_CMD@ src.main:app --host ${EDGE_BIND_HOST} --port ${EDGE_BIND_PORT}
Restart=always
RestartSec=5
StandardOutput=journal
StandardError=journal
LimitNOFILE=65536

[Install]
WantedBy=multi-user.target
